# usbsync64.service
[Unit]
Description=Syncing C64 USB Folder

[Service]
ExecStart=/home/pi/bin/usbsync64

[Install]
WantedBy=default.target

